# clawdis.service
[Unit]
Description=Clawdis Gateway Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/seshu/Documents/Js/clawdis
Environment=PATH=/home/seshu/.nvm/versions/node/v22.21.1/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/pnpm clawdbot gateway --bind lan --force
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
